# datadis-exporter.timer
[Unit]
Description=datadis Exporter Timer

[Timer]
OnCalendar=Weekly
Persistent=true

[Install]
WantedBy=timers.target